# fc-runner.service
[Unit]
Description=fc-runner Firecracker CI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/fc-runner server --config /etc/fc-runner/config.toml
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=300
Environment=RUST_LOG=fc_runner=info

# fc-runner needs root with broad access for:
#   mount/umount (loop devices), chroot (rootfs provisioning),
#   iptables/ipset (NAT rules), TAP devices (ioctl TUNSETIFF)
# Minimal hardening that doesn't break these operations:
ProtectHome=true
ProtectKernelModules=true
RestrictRealtime=true
SupplementaryGroups=kvm

[Install]
WantedBy=multi-user.target
