# fc-runner.service
[Unit]
Description=fc-runner Firecracker CI Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/fc-runner /etc/fc-runner/config.toml
Restart=on-failure
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=300
Environment=RUST_LOG=fc_runner=info

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/fc-runner /opt/fc-runner /tmp
ProtectKernelModules=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictSUIDSGID=true
MemoryDenyWriteExecute=true
# Needed for TAP/NAT and mount/loop operations
# CAP_NET_RAW needed for iptables
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW CAP_SYS_ADMIN CAP_MKNOD
SupplementaryGroups=kvm

[Install]
WantedBy=multi-user.target
